# backend/mockup-backend.service
[Unit]
Description=Mockup Tool Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mockup-tool/backend
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=5
Environment=PORT=4000

[Install]
WantedBy=multi-user.target
